# systemd/lenovo-kernel-updater.timer
[Unit]
Description=Check GitHub Releases for Lenovo CachyOS kernel daily at 4 AM

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true
RandomizedDelaySec=0

[Install]
WantedBy=timers.target

