# systemd/nuc15pro-kernel-updater.timer
[Unit]
Description=Check GitHub Releases for NUC 15 Pro CachyOS kernel daily at 4 AM

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true
RandomizedDelaySec=0

[Install]
WantedBy=timers.target
