# systemd/nuc16pro-servermax-power.service
[Unit]
Description=NUC 16 Pro ServerMax device tuning (BIOS owns power limits and fan curves)
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
# Power limits (PL1/PL2/Tau) and the ACPI platform_profile are deliberately NOT
# written here:
#  - RAPL: the 356H is hard-capped at its 80W Maximum Turbo Power (Intel spec),
#    confirmed on-device (package stays <=80W under stress even with RAPL raised
#    to 104W). A write above 80W is inert and below it would only throttle, so
#    there is nothing to gain; BIOS owns PL1/PL2/Tau.
#  - platform_profile: left to BIOS, which is where the thermal and fan curves
#    are configured on this box. (To bias max turbo from the OS instead, add an
#    ExecStart that writes 'performance' to /sys/firmware/acpi/platform_profile.)
# The updater still reads the RAPL constraints for status, so the BIOS-set
# PL1/PL2 remain visible. Everything below tunes devices WITHIN the BIOS envelope
# and never changes power limits or fan curves.
# energy_perf_bias=0: no microarchitecture power-saving bias on any core
ExecStart=/bin/sh -c 'for b in /sys/devices/system/cpu/cpu*/power/energy_perf_bias; do [ -w "$b" ] && printf 0 > "$b"; done; true'
# NVMe: maximize request queue depth per namespace for Gen4/Gen5 throughput
ExecStart=/bin/sh -c 'for q in /sys/block/nvme*/queue/nr_requests; do [ -w "$q" ] && printf 1023 > "$q"; done; true'
# igc (I226-V 2.5GbE): maximize ring buffers for throughput
ExecStart=/bin/sh -c 'iface=$(ls /sys/class/net/ 2>/dev/null | grep -m1 "^e" || true); [ -n "$iface" ] && ethtool -G "$iface" rx 4096 tx 4096 2>/dev/null || true'
# Thermal: set x86 package passive trip point to 100C = TjMax for Panther Lake 356H.
# This lets the CPU run at full turbo until hardware PROCHOT fires at TjMax.
# Only type=passive trips are touched; type=critical (emergency shutdown) is left untouched.
# Millidegrees: 100000. Writable unconditionally in kernel 6.14+ (no CONFIG_THERMAL_WRITABLE_TRIPS needed).
# NOTE: the loop index is written $${i} so systemd emits a literal ${i} to the shell.
# An unescaped ${i} is expanded by systemd itself to an empty string (journal warns
# "unset environment variable ... i"), which silently breaks the trip-setting loop.
ExecStart=/bin/sh -c 'for zone in /sys/class/thermal/thermal_zone*; do ztype=$(cat "$zone/type" 2>/dev/null || true); [ "$ztype" = "x86_pkg_temp" ] || continue; for i in 0 1; do ttype=$(cat "$zone/trip_point_$${i}_type" 2>/dev/null || true); ttemp="$zone/trip_point_$${i}_temp"; [ "$ttype" = "passive" ] || continue; [ -w "$ttemp" ] && printf 100000 > "$ttemp" 2>/dev/null && echo "thermal: $ttemp=100000" || true; done; done; true'

[Install]
WantedBy=multi-user.target
