# systemd/nuc16pro-kernel-updater.service
[Unit]
Description=Install latest ASUS NUC 16 Pro CachyOS edge kernel from GitHub Releases
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=OWNER_REPO=AmirulAndalib/asus-nuc16pro-cachyos-server-edge-kernel
ExecStart=/usr/local/sbin/nuc16pro-kernel-updater.sh

[Install]
WantedBy=multi-user.target
