# systemd/nuc16pro-servermax-power.service
[Unit]
Description=NUC 16 Pro ServerMax runtime power limits and device tuning
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
# ACPI platform_profile=performance: unified fan/thermal/turbo control on this
# board (firmware/EC owns the actual fan curve; the OS only picks the profile).
# performance holds max turbo residency. Runtime-switchable to balanced/low-power.
ExecStart=/bin/sh -c '[ -f /sys/firmware/acpi/platform_profile ] && echo performance > /sys/firmware/acpi/platform_profile || true'
# Intel RAPL on 120W AC: PL1=95W, PL2=95W, Tau=224s. The Core Ultra 7 356H has a
# Maximum Turbo Power of 80W (Intel spec); on-device the package does not exceed
# ~80W under stress even with RAPL set to 104W, so 80W is the real silicon ceiling.
# RAPL >=95W is therefore intentionally non-binding: it guarantees RAPL never clips
# the CPU below its 80W ceiling, leaving cooling and firmware to govern. Raising
# RAPL past 80W does not add power; the MTP cannot be overridden from the OS.
# BIOS may lock the MSR; read back actual value after write to confirm.
ExecStart=/bin/sh -c 'p=/sys/class/powercap/intel-rapl/intel-rapl:0; [ -d "$p" ] && printf 95000000 > "$p/constraint_0_power_limit_uw" && echo "RAPL PL1=$(cat $p/constraint_0_power_limit_uw)uW" || true'
ExecStart=/bin/sh -c 'p=/sys/class/powercap/intel-rapl/intel-rapl:0; [ -d "$p" ] && printf 95000000 > "$p/constraint_1_power_limit_uw" && echo "RAPL PL2=$(cat $p/constraint_1_power_limit_uw)uW" || true'
ExecStart=/bin/sh -c 'p=/sys/class/powercap/intel-rapl/intel-rapl:0; [ -d "$p" ] && printf 224000000 > "$p/constraint_0_time_window_us" || true'
ExecStart=/bin/sh -c 'p=/sys/class/powercap/intel-rapl/intel-rapl:0; [ -d "$p" ] && printf 224000000 > "$p/constraint_1_time_window_us" || true'
# energy_perf_bias=0: no microarchitecture power-saving bias on any core
ExecStart=/bin/sh -c 'for b in /sys/devices/system/cpu/cpu*/power/energy_perf_bias; do [ -w "$b" ] && printf 0 > "$b"; done; true'
# NVMe: maximize request queue depth per namespace for Gen4/Gen5 throughput
ExecStart=/bin/sh -c 'for q in /sys/block/nvme*/queue/nr_requests; do [ -w "$q" ] && printf 1023 > "$q"; done; true'
# igc (I226-V 2.5GbE): maximize ring buffers for throughput
ExecStart=/bin/sh -c 'iface=$(ls /sys/class/net/ 2>/dev/null | grep -m1 "^e" || true); [ -n "$iface" ] && ethtool -G "$iface" rx 4096 tx 4096 2>/dev/null || true'
# Thermal: set x86 package passive trip point to 100C = TjMax for Panther Lake 356H.
# This lets the CPU run at full turbo until hardware PROCHOT fires at TjMax.
# Only type=passive trips are touched; type=critical (emergency shutdown) is left untouched.
# Millidegrees: 100000. Writable unconditionally in kernel 6.14+ (no CONFIG_THERMAL_WRITABLE_TRIPS needed).
ExecStart=/bin/sh -c 'for zone in /sys/class/thermal/thermal_zone*; do ztype=$(cat "$zone/type" 2>/dev/null || true); [ "$ztype" = "x86_pkg_temp" ] || continue; for i in 0 1; do ttype=$(cat "$zone/trip_point_${i}_type" 2>/dev/null || true); ttemp="$zone/trip_point_${i}_temp"; [ "$ttype" = "passive" ] || continue; [ -w "$ttemp" ] && printf 100000 > "$ttemp" 2>/dev/null && echo "thermal: $ttemp=100000" || true; done; done; true'

[Install]
WantedBy=multi-user.target
